# scripts/build.check.crontab
# run daily at 2am
0 2 * * * $HOME/bin/tokudb.build.check.bash
